# templates/default/carbon-cache.conf.erb
description "Carbon"
author "Chef"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
expect daemon

exec start-stop-daemon --start --oknodo --chdir <%= @home %> --user <%= @user %> --chuid <%= @user %> --pidfile <%= @home %>/storage/carbon-cache-a.pid --name carbon-cache --startas <%= @home %>/bin/carbon-cache.py -- start

emits carbon-cache-running
